# crontab
0 6 * * wed,sat /opt/scripts/geo-ipwhitelist.sh
